# deploy/influence-eth-leaderboards.timer
[Unit]
Description=Update scores at leaderboards

[Timer]
OnBootSec=480s
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target