# deploy/influence-eth-events.timer
[Unit]
Description=Run events fulfill service

[Timer]
OnBootSec=120s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target